# scripts/launchd/com.mssanwel.agent-pm.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mssanwel.agent-pm.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/muhammadsaadshahidanwel/code/agent-pm/scripts/run-weekly.sh</string>
    </array>

    <!-- Monday 09:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/muhammadsaadshahidanwel/code/agent-pm/.claude/agent-pm/logs/launchd-weekly.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/muhammadsaadshahidanwel/code/agent-pm/.claude/agent-pm/logs/launchd-weekly.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/muhammadsaadshahidanwel</string>
    </dict>

    <key>Nice</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
